# file-system/etc/cron.d/wwwsas-post-analyse.example
15 * * * * root "/etc/www-security-assistant/post-analyse.bash" 'AutoMode' >> "/etc/www-security-assistant/tmp/www-security-assistant.exec.log" 2>&1
